# src/systemd/trashcollection.service
[Unit]
Description=Trash Collection Notifier
After=syslog.target network-online.target

[Service]
User=root
ExecStart=/usr/local/bin/trashcollection.jar --spring.config.additional-location=/etc/trashcollection.properties
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
