# application/target/control/deb/postrm
#!/bin/sh

update-rc.d -f thingsboard remove
