# application/target/debian/postinst
#!/bin/sh -e

ec() {
    echo "$@" >&2
    "$@"
}

case "$1" in
    configure)
        

        
        
chown -R thingsboard: /var/log/thingsboard
chown -R thingsboard: /usr/share/thingsboard
update-rc.d thingsboard defaults


        
        ;;
esac
